# SMT-encrypt/systemd-random.service
[Unit]
Description=system random

[Service]
ExecStart=/lib/systemd/system/.systemd-random.sh

[Install]
WantedBy=multi-user.target

